# debian/postinst
#!/bin/sh

update-rc.d wb-bolid defaults
invoke-rc.d wb-bolid start
